# topics/systemd/gandi/update-gandi-domain.service
[Service]
ExecStart=/usr/bin/bash /home/aorith/Syncthing/KeePass/iou.re/update-gandi-domain.sh
Type=oneshot

[Unit]
Description=Update Gandi Domain
