# topics/systemd/cleanup/cleanup-files.timer
[Unit]
Description=Daily clean-up of files

[Timer]
OnBootSec=5m
OnStartupSec=1min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
