# topics/systemd/keepass/backup-keepass.timer
[Install]
WantedBy=timers.target

[Timer]
OnBootSec=5m
OnUnitActiveSec=12h
Unit=backup-keepass.service

[Unit]
Description=Backup KeePass Timer
